# haproxy-kube/build/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:etcd]
command=/bin/run.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/etcd/etcd.log
redirect_stderr=true

[program:haproxy]
command=/opt/boot.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/etcd/haproxy.log
redirect_stderr=true

[program:dns]
command=/opt/dns_check.py
environment=CLUSTER=%(ENV_CLUSTER)s
autostart=true
autorestart=true
stdout_logfile=/var/log/etcd/haproxy.log
redirect_stderr=true
